# local-otel/launchagents/com.frontier.copilot-otel-github-orgs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.frontier.copilot-otel-github-orgs</string>
    <key>ProgramArguments</key>
    <array>
      <string>__LOCAL_OTEL_DIR__/ingest-github-orgs.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true />
    <key>StandardOutPath</key>
    <string>__LOCAL_OTEL_DIR__/launchagent-github-orgs.log</string>
    <key>StandardErrorPath</key>
    <string>__LOCAL_OTEL_DIR__/launchagent-github-orgs.err</string>
  </dict>
</plist>